# demo/deploy/coquic-demo.service
[Unit]
Description=Coquic demo HTTP/3 service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/coquic-demo/current/run-demo.sh
Restart=on-failure
RestartSec=2
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
